# systemd/product-credential-site.service
[Unit]
Description=Product Credential Site
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=www-data
WorkingDirectory=/srv/apps/info.c.bimumedia.com/current
Environment=NODE_ENV=production
EnvironmentFile=/etc/product-credential-site/info.c.bimumedia.com.env
ExecStart=/bin/bash -lc 'exec "$NODE_BIN" "$APP_ROOT/server.js"'
Restart=always
RestartSec=5
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/srv/data/info.c.bimumedia.com
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
